# kalshi-bot.service
[Unit]
Description=Kalshi BTC Mispricing Arbitrage Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/money-money
EnvironmentFile=/home/ubuntu/money-money/.env
ExecStart=/home/ubuntu/money-money/venv/bin/python -m bot.main
Restart=on-failure
RestartSec=30
StandardOutput=append:/home/ubuntu/money-money/logs/bot.log
StandardError=append:/home/ubuntu/money-money/logs/bot.log

[Install]
WantedBy=multi-user.target
